# project/config/production/crontab/map
* * * * * root /bin/bash /var/www/map/project/tool/production/check_update.sh
